# systemd/so-health.timer
[Unit]
Description=Run SO daily health report at 7:10 AM EST

[Timer]
OnCalendar=*-*-* 12:10:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
